# .docker/config/supervisor.d/cron.ini
[program:cron]
command = crond -f
autostart=%(ENV_ENABLE_CRON_JOB)s
autorestart=true
