# .docker/config/supervisor.d/horizon.ini
[program:horizon]
process_name=%(program_name)s
command=php /app/artisan horizon
autostart=%(ENV_ENABLE_LARAVEL_HORIZON)s
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/app/storage/horizon.log
