# .docker/config/supervisor.d/laravel-worker.ini
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work
autostart=%(ENV_ENABLE_LARAVEL_WORKER)s
autorestart=true
user=www-data
numprocs=8
redirect_stderr=true
stdout_logfile=/app/storage/logs/worker.log
